# nrf24-timed.service
[Unit]
Description=Time server publishing current date/time via NRF24L01.

